# profile.d/fzf.sh
#!/bin/bash

source /etc/fzf/.fzf.config
